# scripts/host-watchdog/parley-fs-watchdog.timer
[Unit]
Description=Run parley-fs-watchdog every 30 seconds

[Timer]
OnBootSec=2min
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
