# scripts/host-watchdog/parley-fs-watchdog.service
[Unit]
Description=Parley filesystem watchdog (hard-reboot on EIO/EROFS)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/parley-fs-watchdog
